# agent/contrib/systemd/opensource-agent.service
[Unit]
Description=MIE opensource-server agent check-in
Wants=network-online.target
# environment.service (container images) writes the runtime env into
# /etc/environment; harmless ordering hint when the unit doesn't exist.
After=network-online.target environment.service

[Service]
Type=oneshot
# Runtime config (SITE_ID, MANAGER_URL, API_KEY) lives in /etc/environment.
EnvironmentFile=-/etc/environment
ExecStart=/usr/bin/node /opt/opensource-server/agent/dist/index.js
