# images/base/sssd-autoconfigure.service
[Unit]
Description=Auto-configure SSSD from DHCP search domain
After=network-online.target
Before=sssd.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/sssd-autoconfigure

[Install]
WantedBy=multi-user.target
