# nginx-reverse-proxy/pull-config.cron
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
* * * * * root /opt/opensource-server/nginx-reverse-proxy/pull-config.sh
